# SDK/Belkin_Plug-Ins_20111115_v0.00.10/user/sysinit/inittab
# slog:unknown:/sbin/syslogd -n
# klog:unknown:/sbin/klogd -n
